# files/rc.local
#!/bin/sh

# sudo systemctl start rc-local.service

/usr/local/bin/uwsgi --emperor /etc/uwsgi/vassals/ &
